# debian/pki-kra.postinst
#!/bin/sh
set -e

invoke_failure() {
    # invoke-rc.d failed, likely because no instance has been configured yet
    # but exit with an error if an instance is configured and the invoke failed
    if [ ! -d /etc/pki-kra ]; then
        echo "... because no KRA instance has been configured yet."
    else
        exit 1
    fi
}

#DEBHELPER#
